# deploy/ip-https/mock-interview-cert-renew.cron
17 */12 * * * root /opt/mock-interview-proxy/renew-ip-certificate.sh >>/var/log/mock-interview-cert-renew.log 2>&1
